# etc/supervisor/conf.d/fb.conf
[program:fb]
command = /bin/bash /opt/fb/env.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/fb.err.log
stdout_logfile=/var/log/fb.out.log
